# backup/080911/config/cron.hourly/ispmgr_rotate
#!/bin/sh

/usr/local/ispmgr/sbin/rotate

